# docker/root/etc/services.d/traefik-certificate-extractor/run
#!/usr/bin/with-contenv bash

cd /app || exit

exec \
    s6-setuidgid abc /app/run.sh